# deploy/esv-deploy.service
# Installed at /etc/systemd/system/esv-deploy.service on the VPS.
# Runs the pull-based deploy poller once (triggered by esv-deploy.timer).
# Requires (one-time host setup):
#   git config --system --add safe.directory /home/englishspeaking/englishspeakinggermany.online
[Unit]
Description=Pull-based deploy for englishspeakinggermany.online
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
# Least privilege: npm ci/build lifecycle scripts must NOT run as root.
# The englishspeaking user owns the web root; systemctl restart vet-api is
# delegated via a narrow sudoers rule (see deploy/README or DECISIONS.md):
#   englishspeaking ALL=(root) NOPASSWD: /usr/bin/systemctl restart vet-api
User=englishspeaking
Group=englishspeaking
Environment=HOME=/home/englishspeaking
ExecStart=/usr/local/bin/esv-poll-deploy.sh
