# api/vet-api.service
[Unit]
Description=Gunicorn instance to serve The Pack API
After=network.target

[Service]
# Least privilege: run as dedicated app user, not root
User=englishspeaking
Group=www-data
WorkingDirectory=/home/englishspeaking/englishspeakinggermany.online/api
Environment="PATH=/home/englishspeaking/englishspeakinggermany.online/api/venv/bin"
Environment="PORT=5000"
EnvironmentFile=/home/englishspeaking/englishspeakinggermany.online/api/.env

ExecStart=/home/englishspeaking/englishspeakinggermany.online/api/venv/bin/gunicorn --workers 3 --bind 127.0.0.1:5000 main:app

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/home/englishspeaking/englishspeakinggermany.online/api
PrivateTmp=true

[Install]
WantedBy=multi-user.target
